# deploy/flask-starter.service
[Unit]
Description=Flask Starter App
After=network.target

[Service]
User=<your-linux-user>
WorkingDirectory=/path/to/flask-starter
EnvironmentFile=/path/to/flask-starter/.env
ExecStart=/path/to/flask-starter/.venv/bin/gunicorn \
    -w 2 \
    -b <tailscale-ip>:<port> \
    "app:create_app()"
Restart=on-failure

[Install]
WantedBy=multi-user.target
